# deploy/sub2api-maintenance.cron
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

*/5 * * * * root /usr/local/sbin/sub2api-healthcheck >/dev/null 2>&1
17 4 * * 0 root /usr/bin/docker image prune --force --filter until=168h >/dev/null 2>&1
